# X11/dot-local/share/X11/sv/.s6-svscan/finish
#!/bin/sh

[ -n "$X11_SV_DIR" ] && rm -rf "$X11_SV_DIR"
